# server/sensor_hub.service
[Unit]
Description=Sensor Hub API server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/sensor_projekt_V2/server
ExecStart=/home/pi/sensor_projekt_V2/venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5

# Optional overrides — uncomment and edit as needed:
# Environment=DB_PATH=/home/pi/sensor_hub.db
# Environment=POLL_INTERVAL_S=1800
# Environment=DASHBOARD_DIR=/home/pi/sensor_projekt_V2/dashboard

[Install]
WantedBy=multi-user.target
